# service-scripts/bbb-timing.service
# This file goes into /etc/systmed/system/
[Unit]
Description=Arbwave BBB backend for the Timing device
After=network.target network-online.target auditd.service
Requires=afrl-timing-fw.service
After=afrl-timing-fw.service

[Service]
EnvironmentFile=/etc/default/arbwave
User=arbwave
WorkingDirectory=/tmp/
# give up to 10 seconds for PRU to boot
ExecStartPre=/bin/sh -c \
  "for i in 1 2 3 4 5; do if [ -c /dev/rpmsg_pru30 ] ; then break; else sleep 5; fi ; done"
ExecStartPre=/usr/bin/test -e /dev/rpmsg_pru30
ExecStart=/usr/bin/arbwave-bbb_timing_controller $ARBWAVE_BBB_CNTRL_ARGS
TimeoutStopSec=8
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
Alias=bbb-timing.service
